# deployment/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid

[program:php-fpm]
command=/bin/bash -c "/usr/sbin/php-fpm7.2 -F"
